# systemd/telegram-selling-bot.service.example
[Unit]
Description=Telegram Selling Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USER
WorkingDirectory=/opt/telegram-selling-bot
EnvironmentFile=/opt/telegram-selling-bot/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/telegram-selling-bot/.venv/bin/python /opt/telegram-selling-bot/selling_bot.py
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
